# kiosk_files/busstop.service
[Unit]
Description=busstop
After=network.target

[Service]
User=pi
WorkingDirectory=/home/pi/busstop/flask_app
ExecStart=/home/pi/busstop/flask_app/venv/bin/gunicorn --preload --bind localhost:8000 --workers 1 app:app
Restart=always

[Install]
WantedBy=multi-user.target
